# root/etc/services.d/haproxyctl/run
#!/usr/bin/with-contenv sh

/var/omega-haproxyctl/omega-haproxyctl -bind="${BIND:-:5004}" -config=${CONFIG_PATH:-/config/production.json}
